# deployment/supervisor/blog.ini
[program:blog]
directory=/home/alex/flask_blog                     ; Location of application
command=/home/alex/.local/share/virtualenvs/flask_blog-Zgg6zSMC/bin/gunicorn -b localhost:5000 --log-level debug --timeout 90 --error-logfile /home/alex/log/gun3.log -k gevent -w 2 --threads 100 g_server:app                   ; The command to execute
autostart=true                           ; Start this application when supervisord starts
autorestart=true                         ; Restart this application if it crashes
stderr_logfile=/home/alex/log/blog.err.log  ; Make sure this directory exists
stdout_logfile=/home/alex/log/blog.log      ; Make sure this directory exists
stopsignal=INT                           ; Signal sent to the application when halting
user=alex                           ; setuid to this UNIX account to run the program

